# usr/lib/systemd/system/fogbot.service
[Unit]
Description=fogbot - Linux intrusion detection via Telegram
Documentation=https://github.com/penguinpowernz/fogbot
After=network-online.target auditd.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fogbot daemon
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/fogbot /var/log/fogbot

# Capabilities required for monitoring
# Adjust based on enabled skills
AmbientCapabilities=CAP_NET_ADMIN CAP_AUDIT_CONTROL CAP_SYS_PTRACE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_AUDIT_CONTROL CAP_SYS_PTRACE

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fogbot

# Reload on SIGHUP
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
